# config/systemd/user/usb-remove.service
[Unit]
Description=Play USB sound

[Service]
Type=oneshot
ExecStart=/usr/bin/paplay /usr/share/sounds/usb-disconnect.wav
